# debezium-postgres-citus/xinetd.d/pgsqlchk
service pgsqlchk
{
        flags           = IPv4
        socket_type     = stream
        port            = 9201
        wait            = no
        user            = postgres
        server          = /opt/pgsqlchk
        log_on_failure  += USERID
        disable         = no
        only_from       = 0.0.0.0/0
        per_source      = UNLIMITED
}